# states/fm/player/files/soundwave.init.conf
#
# Upstart Script for FM SoundWave
#

description "SOON_ FM SoundWave"

start on (started networking)

env PATH=$PATH:/soundwave/bin

respawn

script
    exec soundwave
end script
